# 50_Use_AWS_Environment/roles/spring-boot/templates/spring-boot.service.j2
[Unit]
Description={{app.desc}}
After=syslog.target

[Service]
User={{app.user}}
ExecStart=/home/{{app.user}}/{{app.filename}} SuccessExitStatus=143

[Install]
WantedBy=multi-user.target